# supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/sbin/service mysqld start

[program:httpd]
command=/sbin/service httpd start
